# noninteractive/.crontab
MAILTO=
#
#mmhh DD MM d (mm=minute, hh=hour, DD=DayOfMonth, MM=MonthNum, d=DayOfWeekWhereZeroIsSunday)
@reboot /cygdrive/c/Users/KevinAndrewLipscomb/Documents/SANDBOX/vocational/kalipso-infogistics/OWNREPO/Oscalerter/bin/Oscalerter.exe
