# noninteractive/.crontab
MAILTO=
#
#mmhh DD MM d (mm=minute, hh=hour, DD=DayOfMonth, MM=MonthNum, d=DayOfWeekWhereZeroIsSunday)
@reboot /cygdrive/c/PROJ/Oscalerter/bin/Debug/Oscalerter.exe
